# build/z-way/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:z-way-server]
command=/opt/z-way-server/z-way-server --log-level 2 --log-file /var/log/z-way-server/logfile
directory=/opt/z-way-server
autorestart=true
autostart=true

# Redirect error to supervisord
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:z-way-server-log]
command=tail -f /var/log/z-way-server/logfile
autorestart=true
autostart=true

# Redirect output to supervisord
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

# Redirect error to supervisord
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0